# nextcloud/crontab
*/5 * * * * php -f /var/www/html/cron.php
40 19 * * * bash -c "PHP_MEMORY_LIMIT=512M php /var/www/html/occ preview:pre-generate -vvv"