# nextcloud/crontab
*/5 * * * * php -f /var/www/html/cron.php
30 19 * * * php /var/www/html/occ preview:pre-generate -vvv