# crontab/crontab
* * * * * echo "Crontab is working" > /proc/1/fd/1
